# root/etc/services.d/xteve/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc "${XTEVE_BIN}/xteve" \
	-config=/config \
	-port="$XTEVE_PORT"
