# root/etc/services.d/zap2xml/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /zap2xml/zap2xml
	
