# root/etc/services.d/zap2xml/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc $ZAP2XML_HOME/zap2xml
	